# contrib/ement-e2ee.service
# systemd user unit: cp to ~/.config/systemd/user/ then
#   systemctl --user enable --now ement-e2ee
# On WSL2, enable systemd in /etc/wsl.conf ([boot] systemd=true) or run the
# proxy from your shell profile instead.
[Unit]
Description=Transparent E2EE proxy for ement.el
After=network-online.target

[Service]
ExecStart=%h/.cargo/bin/ement-e2ee serve
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
